# Linux/TP3/3/db_dump.timer
[Unit]
Description=Run service db_dump

[Timer]
OnCalendar=*-*-* 4:00:00

[Install]
WantedBy=timers.target